# .config/systemd/user/cec.service
[Unit]
Description=CEC initialisation
After=systemd-user-sessions.service

[Install]
WantedBy=default.target

[Service]
Type=simple
ExecStart=/usr/bin/cec-ctl -o RPi --playback -S
Restart=on-failure
RestartSec=10
TimeoutStopSec=10
