# .config/pp_client_download.timer
[Unit]
Description=Schedule PhotoPrism Download

[Timer]
OnCalendar=*-*-* 00:01:00

[Install]
WantedBy=default.target
